# config/systemd/user/my-org-export.service
[Unit]
Description=Exports dotfiles org

[Service]
ExecStart=/home/green/scripts/emacs-publish-dotfiles.sh
